# forum-learning/nginx/crontab
* * * * * ./certificate/generate_cert.sh
